# kafka-broker/supervisor/node_exporter.conf
[program:node_exporter]
command=/usr/bin/node_exporter
autostart=true
autorestart=true
